# deploy/systemd/damian-www-restart.service
[Unit]
Description=Daily restart of damian_www (clears any slow leaks, picks up new config)

[Service]
Type=oneshot
ExecStart=/bin/systemctl try-restart damian-www.service
